# backend/systemd/artframe.service
[Unit]
Description=Artframe Digital Photo Frame
Documentation=https://github.com/yourusername/artframe
After=network-online.target
Wants=network-online.target
StartLimitBurst=5
StartLimitIntervalSec=600

[Service]
Type=simple
User=pi
Group=pi
# Use data dir as working directory so lgpio can create notification pipes
WorkingDirectory=/opt/artframe/backend/data

# Environment - include both possible uv install locations
Environment="PATH=/home/pi/.local/bin:/home/pi/.cargo/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONUNBUFFERED=1"
# Use system cache dir since ProtectHome=read-only prevents writing to ~/.cache
Environment="UV_CACHE_DIR=/var/cache/artframe/uv"
# Force lgpio as the GPIO pin factory (prevents fallback to incompatible factories)
Environment="GPIOZERO_PIN_FACTORY=lgpio"

# Main process - using uv run for dependency management
# Use bash -c to resolve uv from PATH (supports both ~/.local/bin and ~/.cargo/bin)
ExecStart=/bin/bash -c 'exec uv run artframe --config /opt/artframe/backend/config/artframe.yaml'
# Display cleanup is handled in the main process shutdown (FastAPI lifespan)
# ExecStopPost is unreliable due to vendor epdconfig claiming GPIO at import time

# Restart policy
Restart=always
RestartSec=10

# Graceful shutdown
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

# Allow binding to port 80 and raw hardware I/O for SPI/GPIO
AmbientCapabilities=CAP_NET_BIND_SERVICE CAP_SYS_RAWIO

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/log/artframe /var/cache/artframe /opt/artframe/backend/data

# GPIO and SPI access for e-ink display
SupplementaryGroups=gpio spi
DevicePolicy=auto

# Resource limits
MemoryMax=512M
CPUQuota=50%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=artframe

[Install]
WantedBy=multi-user.target
