# backend/systemd/artframe.service
[Unit]
Description=Artframe Digital Photo Frame
Documentation=https://github.com/yourusername/artframe
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/artframe/backend

# Environment - include both possible uv install locations
Environment="PATH=/home/pi/.local/bin:/home/pi/.cargo/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONUNBUFFERED=1"

# Main process - using uv run for dependency management
ExecStart=/home/pi/.local/bin/uv run artframe --config /opt/artframe/backend/config/artframe.yaml
# Use "-" prefix so service stop doesn't fail if display clear fails
ExecStopPost=-/home/pi/.local/bin/uv run python -m artframe.clear_display

# Restart policy
Restart=always
RestartSec=10
StartLimitBurst=5
StartLimitIntervalSec=600

# Graceful shutdown
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

# Allow binding to port 80
AmbientCapabilities=CAP_NET_BIND_SERVICE

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/log/artframe /var/cache/artframe /opt/artframe/backend/data

# Resource limits
MemoryMax=512M
CPUQuota=50%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=artframe

[Install]
WantedBy=multi-user.target
